# scripts/systemd/units/futon-services.slice
[Unit]
Description=Protected Futon serving workloads

[Slice]
CPUAccounting=yes
MemoryAccounting=yes
IOAccounting=yes
TasksAccounting=yes
ManagedOOMPreference=avoid
